# bin/phoenix.hal9000.service
[Unit]
Description=Phoenix Monitor Flask Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/stefano/Documents/MySkynet/src/PhoeNix/server
ExecStart=/home/stefano/Documents/MySkynet/src/PhoeNix/server/phoenixenv/bin/python phoenix_server.py
Restart=always
User=stefano

[Install]
WantedBy=multi-user.target
